# config/zoom-transcript/com.zoom-transcript.mover.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zoom-transcript.mover</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-l</string>
    <string>{{ ansible_env.HOME }}/.dotfiles/scripts/move-zoom-transcripts.sh</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>{{ ansible_env.HOME }}/Downloads</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>/tmp/zoom-mover.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/zoom-mover.err.log</string>
</dict>
</plist>
